# services/eztv-rss-to-json/unit/eztv.timer
[Unit]
Description=EZTV RSS to JSON

[Timer]
OnCalendar=*:0/1
Unit=eztv-rss.service

[Install]
WantedBy=timers.target