# units/summer-pockets.service
[Unit]
Description=Summer Pockets Auto Downloader
After=network.target

[Service]
Type=oneshot
User=keeb
Group=keeb
WorkingDirectory=/storage/02/code/projects/media-management-service
Environment=PYTHONPATH=/storage/02/code/projects/media-management-service
ExecStart=/usr/bin/python3 services/subsplease-show-page-to-json/show-1080p-magnet.py summer-pockets --auto
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target